# linux/systemd/.config/systemd/user/mpris-proxy.service
[Unit]
Description=Forward bluetooth media controls to MPRIS

[Service]
Type=simple
ExecStart=/sbin/mpris-proxy

[Install]
WantedBy=default.target
