# linux/systemd/.config/systemd/user/notes-sync.service
[Unit]
Description=Run git sync on notes

[Service]
Type=oneshot
ExecStart=/home/laznp/.local/bin/notes-sync

