# linux/systemd/.config/systemd/user/polkit-kde-authentication-agent.service
[Unit]
Description=Polkit KDE Authentication Agent
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
ExecStart=/usr/lib/polkit-kde-authentication-agent-1
Restart=on-failure
RestartSec=2
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority

[Install]
WantedBy=default.target
