# linux/systemd/.config/systemd/user/udiskie.service
[Unit]
Description=udiskie — automount removable media

[Service]
ExecStart=/usr/sbin/udiskie
Restart=on-failure

[Install]
WantedBy=default.target